# recipes-core/udev/files/99-toradex-usb-udc.rules
ACTION=="add", SUBSYSTEM=="udc", RUN+="/usr/bin/gadget-start"
